# prometheus-bmp180-exporter.service
[Unit]
Description=Prometheus BMP180 exporter
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python /usr/bin/prometheus-bmp180-exporter.py
ProtectHome=true
NoNewPrivileges=true
RestrictNamespaces=yes
RestrictRealtime=yes
ProtectKernelModules=true
ProtectKernelTunables=true
LockPersonality=true

[Install]
WantedBy=multi-user.target
